# debian/templates/rc.local.erb
#!/bin/sh -e

if [ -x "/usr/sbin/logwatch" ]; then
  /usr/sbin/logwatch --output mail
fi

exit 0
