# debian/templates/rc.local.erb
#!/bin/sh -e

if [ -x "/etc/cron.daily/00logwatch" ]; then
  /etc/cron.daily/00logwatch
fi

if [ -x "/etc/cron.daily/rkhunter" ]; then
  /etc/cron.daily/rkhunter
fi

exit 0
